# packaging/linux/systemd/openlogi-agent.service
# OpenLogi background agent — packaged user unit.
#
# Packages use /usr/bin directly. install.sh rewrites ExecStart when installing
# under a custom PREFIX (default: /usr/local).
# The agent's launch_at_login setting writes its own copy to
# $XDG_CONFIG_HOME/systemd/user/ using the current binary path — both coexist,
# enabling either works.
#
# Enable for the current user:
#   systemctl --user enable --now openlogi-agent.service

[Unit]
Description=OpenLogi background agent (Logitech HID++ device control)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/openlogi-agent
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
